# sys_files/usr/lib/systemd/system/llama-server.service
[Unit]
Description=llama.cpp OpenAI-compatible API server
After=network-online.target
Wants=network-online.target
ConditionPathExists=/var/lib/prism/models/default.gguf

[Service]
Type=simple
EnvironmentFile=/usr/lib/prism/llama-server.env
WorkingDirectory=/usr/lib/llama.cpp
ExecStart=/usr/lib/llama.cpp/llama-server \
    --host 127.0.0.1 \
    --port 8080 \
    --model /var/lib/prism/models/default.gguf \
    --ctx-size ${PRISM_CTX_SIZE} \
    --n-gpu-layers ${PRISM_GPU_LAYERS}
Restart=on-failure

[Install]
WantedBy=multi-user.target
