# sys_files/usr/lib/systemd/system/llama-server.path
[Unit]
Description=Start llama-server when the default model is available

[Path]
PathExists=/var/lib/llama/models/default.gguf

[Install]
WantedBy=multi-user.target
