# sys_files/etc/xdg/autostart/prism-llama-model.desktop
[Desktop Entry]
Type=Application
Name=Prism AI Model Setup
Comment=Download local AI model on first use
Exec=/usr/libexec/prism/ensure-llama-model --gui
Terminal=false
NoDisplay=true
X-GNOME-Autostart-enabled=true
